# tty/.config/systemd/user/hermes.service
[Unit]
Description=Hermes refresh info and MOTD

[Service]
Type=oneshot
RemainAfterExit=true
Environment="TTY_DOTFILES=%h/.orpheus/projects/dotfiles"
ExecStop=${TTY_DOTFILES}/tty/.local/bin/hermes
StandardOutput=truncate:${TTY_DOTFILES}/tty/.config/hermes/motd.log

[Install]
WantedBy=default.target
