# tty/.config/systemd/user/restic-check-read-local.timer
[Unit]
Description=Monthly check read for restic local repository

[Timer]
OnCalendar=*-*-09 19:00:00
# RandomizedDelaySec=1h
Persistent=true
Unit=restic-check-read@local.service

[Install]
WantedBy=timers.target
