# tty/.config/systemd/user/sway.target
[Unit]
Description=Sway compositor session
Documentation=man:systemd.special
BindsTo=graphical-session.target
Wants=graphical-session.target
After=graphical-session.target
