# tty/.config/systemd/user/switcher.service
[Unit]
Description=Switcher theme mode next session

[Service]
Type=oneshot
RemainAfterExit=true
# Environment variables are listed here as a fallback
# documentation. They are imported via `systemctl --user
# import-environment` in `.bash_profile` and Sway config.
EnvironmentFile=%h/.config/systemd/envs/switcher.env
ExecStart=%h/.local/bin/switcher_helper

[Install]
WantedBy=default.target
